# dell-fan.timer
[Unit]
Description=Timer para manter Fan Boost e Perfil Custom

[Timer]
OnBootSec=5
OnUnitActiveSec=5

[Install]
WantedBy=timers.target
